# twitter2stream.conf
exec python /home/ubuntu/twitter2stream.py

start on filesystem or runlevel [2345] or stopped twitter2stream
stop on shutdown
